# clients/deploy/rpi4/omega-voice.service
[Unit]
Description=OmegaAgent RPi 4 Voice Assistant
After=pulseaudio.service network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/jaeminbbq/Projects/OmegaAgent/clients
EnvironmentFile=/home/jaeminbbq/Projects/OmegaAgent/clients/.env
ExecStart=/home/jaeminbbq/Projects/OmegaAgent/clients/.venv/bin/python3 rpi4_assistant.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
